# deployments/systemd/vpanel-agent.service
[Unit]
Description=V Panel Node Agent
Documentation=https://github.com/your-repo/V
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vpanel-agent
ExecStart=/opt/vpanel-agent/vpanel-agent -config /opt/vpanel-agent/configs/agent.yaml
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/vpanel-agent /etc/xray /var/backups/xray

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
